# services/ingestion/systemd/newsroom-poller.timer
[Unit]
Description=Poll Google News every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=30s
Unit=newsroom-poller.service

[Install]
WantedBy=timers.target
